# data/docker/web/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
autorestart=true
autostart=true
killasgroup=true
stopasgroup=true

[program:cron]
command=cron -f -L 15
autostart=true
autorestart=true

[program:queue-work]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --tries=1
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue-work.log
